# pmis_v2/deploy/launchd/com.pmis.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pmis.sync</string>

  <key>ProgramArguments</key>
  <array>
    <string>__REPO_ROOT__/pmis_v2/deploy/launchd/pmis_sync.sh</string>
  </array>

  <key>StartInterval</key>
  <integer>1800</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__/pmis_v2</string>

  <key>StandardOutPath</key>
  <string>__REPO_ROOT__/pmis_v2/logs/scheduler/sync.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>__REPO_ROOT__/pmis_v2/logs/scheduler/sync.stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
