# pmis_v2/deploy/launchd/com.pmis.sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pmis.sync</string>

  <!-- Call python3 directly — shell invocation hits TCC "Operation not
       permitted" when the script lives under ~/Desktop. Python binaries
       on macOS inherit Full Disk Access from the user's PATH-granted
       interpreters, matching the existing com.pmis.memory-server pattern. -->
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/python3</string>
    <string>__REPO_ROOT__/pmis_v2/cli.py</string>
    <string>sync</string>
    <string>run</string>
    <string>--model</string>
    <string>qwen2.5:7b</string>
  </array>

  <key>StartInterval</key>
  <integer>1800</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>WorkingDirectory</key>
  <string>__REPO_ROOT__/pmis_v2</string>

  <key>StandardOutPath</key>
  <string>__REPO_ROOT__/pmis_v2/logs/scheduler/sync.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>__REPO_ROOT__/pmis_v2/logs/scheduler/sync.stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>PYTHONPATH</key>
    <string>__REPO_ROOT__/pmis_v2</string>
  </dict>
</dict>
</plist>
