# pmis_v2/desktop/com.pmis.menubar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pmis.menubar</string>

  <key>ProgramArguments</key>
  <array>
    <string>__VENV_PYTHON__</string>
    <string>__SCRIPT_PATH__</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
  </dict>

  <!-- Rate limit restart attempts so a broken install doesn't spin the CPU -->
  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>__LOG_DIR__/menubar.out.log</string>
  <key>StandardErrorPath</key>
  <string>__LOG_DIR__/menubar.err.log</string>

  <!-- GUI apps need this to attach to the user session; otherwise the icon
       won't show up in the menu bar. -->
  <key>ProcessType</key>
  <string>Interactive</string>

  <key>LimitLoadToSessionType</key>
  <string>Aqua</string>
</dict>
</plist>
